# etc/local_bin/profile.d/local_bin.sh
# extending $PATH to ~/.local/bin, for Claude!
export PATH="$HOME/.local/bin:$PATH"
